# rootfs/etc/s6-overlay/s6-rc.d/web-ui/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Web UI Server
# Serves the dashboard on port 8099; proxies /api/* to agent API on 8098
# ==============================================================================

bashio::log.info "Starting Web UI on port 8099..."

cd /app

exec python3 -m src.web_ui_server
