# rootfs/etc/services.d/scan-agent/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Scan Agent finish script
# Called when service stops
# ==============================================================================

s6-echo "[finish] Scan Agent stopped"
