# rootfs/etc/services.d/web-ui/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Web UI finish script
# Called when service stops
# ==============================================================================

s6-echo "[finish] Web UI stopped"
