# rootfs/etc/s6-overlay/s6-rc.d/init-prepare/finish
#!/usr/bin/execlineb -S0
s6-echo "[finish] init-prepare done"
