# rootfs/etc/services.d/ftp-server/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# FTP Server finish script
# Called when service stops
# ==============================================================================

s6-echo "[finish] FTP server stopped"
