# rootfs/etc/services.d/ftp-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start FTP Server
# Credentials loaded from S6 container environment (set by 00-prepare.sh)
# Auto-restarts if crashes
# ==============================================================================

bashio::log.info "Starting FTP server on port 2121..."

cd /app

exec python3 -m src.agent.ftp_server
