# rootfs/etc/s6-overlay/s6-rc.d/scan-agent/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Scan Agent
# Config is at /data/config.yaml — generated from HA options by init-prepare
# Exposes internal API on 127.0.0.1:8098
# Auto-restarts if crashes
# ==============================================================================

bashio::log.info "Starting Scan Agent..."

cd /app

exec python3 -m src.main --config /data/config.yaml
